# debian/creature-mqtt.service
[Unit]
Description=Creature MQTT Bridge
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/default/creature-mqtt
ExecStart=/usr/bin/creature-mqtt bridge \
    --host ${CREATURE_SERVER_HOST} \
    --port ${CREATURE_SERVER_PORT} \
    --mqtt-host ${MQTT_HOST} \
    --mqtt-port ${MQTT_PORT}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
